# supervisor.conf
[program:hubot]
environment=HUBOT_SLACK_TOKEN={{ slack_token }}
directory={{ app_path }}
command={{ app_path }}/bin/hubot --adapter slack
startsecs = 0
